# qnx/wm_inst/etc/systemd/system/getty@tty1.service.d/99-weldmaster.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --nohostname --noissue --autologin precitec %I $TERM
