# qnx/wm_inst/etc/systemd/system/cpugov-performance.service.in
[Unit]
Description=Switches CPU governor to performance

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@CMAKE_INSTALL_PREFIX@/opt/wm_inst/batch/cpugov-performance.sh

[Install]
WantedBy=multi-user.target

