# qnx/wm_inst/etc/systemd/system/weldmaster-network-devices.service.in
[Unit]
Description=Configure Weldmaster Network Devices on system startup

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@CMAKE_INSTALL_PREFIX@/opt/configure-network-devices.sh

[Install]
WantedBy=multi-user.target
